# MUX-Nantes/config/supervisor/enc-radio8.conf
[program:enc-radio8]
command=/usr/local/bin/odr-audioenc -v http://streaming.prun.net/streamprun -r 48000 -c 2 -o tcp://localhost:9008 -b 128 -w radio8.dls -p 34 -P radio8.pad slide -s 60
directory=/home/odr/config/mot/radio8
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/enc-radio8.log
stdout_logfile=/var/log/supervisor/enc-radio8.log

[program:mot-radio8]
command=odr-padenc -o radio8.pad -t radio8.dls -c 15 -p 34 -d slide
directory=/home/odr/config/mot/radio8
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/mot-radio8.log
stdout_logfile=/var/log/supervisor/mot-radio8.log
